# private_dot_config/systemd/user/battery-status.service
[Unit]
Description=Show battery percentage

[Service]
Type=oneshot
ExecStart=%h/scripts/notify-battery-status.sh
